# deploy/arxiv-rag-frontend.service
[Unit]
Description=arXiv RAG Next.js frontend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/ChatExamenRI
ExecStart=/usr/bin/npm run start -- --port 3000
Restart=on-failure
RestartSec=5
User=ubuntu
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
