# deploy/arxiv-rag-backend.service
[Unit]
Description=arXiv RAG FastAPI backend
After=network.target

[Service]
Type=simple
# Ajusta esta ruta si clonaste el repo en otro directorio del usuario ubuntu.
WorkingDirectory=/home/ubuntu/ChatExamenRI
Environment=GEMINI_API_KEY=REEMPLAZA_CON_TU_CLAVE_REAL
Environment=CORS_ORIGINS=*
ExecStart=/home/ubuntu/ChatExamenRI/backend/.venv/bin/python -m uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
User=ubuntu

[Install]
WantedBy=multi-user.target
